# java-addon/rootfs/etc/services.d/java-run/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the clever-home service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================


CONFIG_DIR="/config"
APPS_DIR="$CONFIG_DIR/apps"
mkdir -p "$APPS_DIR"

# Copiază JAR-ul default doar dacă nu există deja
if [ ! -f "$APPS_DIR/app.jar" ]; then
    cp /opt/app/app.jar "$APPS_DIR/app.jar"
fi



# Declare variables
declare ha_host
declare ha_token

## Get the config options
ha_host=$(bashio::config 'ha_host')
ha_token=$(bashio::config 'ha_token')


# Export Home Assistant variables for Java
export HA_HOST="$ha_host"
export HA_TOKEN="$ha_token"

## Run Java program
exec java -jar "$APPS_DIR/app.jar"